# startup/FIRSTBOOT.service
[Unit]
Description=FIRSTBOOT
After=default.target
[Service]
ExecStart=/root/FIRSTBOOT.sh
[Install]
WantedBy=default.target